# meridian.service
[Unit]
Description=MERIDIAN ML Pipeline Service
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/dev/meridian
Environment="PATH=/home/%i/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/%i/dev/meridian"
EnvironmentFile=-/home/%i/dev/meridian/.env

# Use Docker Compose for full stack
ExecStart=/usr/bin/docker compose -f /home/%i/dev/meridian/deploy/docker-compose.dgx.yml up
ExecStop=/usr/bin/docker compose -f /home/%i/dev/meridian/deploy/docker-compose.dgx.yml down

# Restart policy
Restart=on-failure
RestartSec=10s
TimeoutStartSec=300

# Resource limits (adjust for DGX capabilities)
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target